# certbot/cert-relay.service
[Unit]
Description=Certificate ACL Relay
Documentation=Reads domain names from a FIFO and applies filesystem ACLs for renewed certs
After=network.target

[Service]
Type=simple
ExecStart=/opt/certbot/cert-relay-host.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
